# services/code-server/finish
#!/data/data/com.termux/files/usr/bin/sh
pkill code-server
