# services/code-server/run
#!/data/data/com.termux/files/usr/bin/sh

code-server --config ./config.yml
