# services/code-server/run
#!/data/data/com.termux/files/usr/bin/sh

code-server --bind-addr 0.0.0.0:8080 --disable-telemetry
