# services/dockerd/finish
#!/data/data/com.termux/files/usr/bin/sh
su -c pkill dockerd
